# LEGACY_IMPORT/Legacy NeuronX/mac-setup/com.ollama.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ollama.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ranjansingh/Desktop/NeuronX/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ranjansingh/Desktop/NeuronX</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/ranjansingh/Desktop/NeuronX</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ollama-gateway.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ollama-gateway.err</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StartInterval</key>
    <integer>300</integer>
</dict>
</plist>